# docker/crawler/crontab
0 7  * * * /app/docker/crawler/run-crawl.sh
30 15 * * * /app/docker/crawler/run-crawl.sh
